# module/tools/SystemConfiguration/data/system/default/etc/systemd/system/robocup.service.d/usb_no_restart.conf
# This adds a condition to the robocup service to prevent it from starting if
# any USB devices are connected that are not in the allowlist. This is useful
# for preventing the service from restarting mid-competition while we are
# fixing a robot - any keyboard, usb monitor, etc. will stop it.

[Unit]
Requires=robocup-usb-check.service
After=robocup-usb-check.service
